# package/debian/DEBIAN/postinst
#!/bin/sh
set -e

SERVICE="smartdns"
# 读取 preinst 保存的状态文件（与 preinst 中定义的路径一致）
STATE_FILE="/var/run/smartdns-upgrade-state"

# 重新加载 systemd 配置（确保新服务文件生效）
systemctl daemon-reload >/dev/null 2>&1 || true

# 处理升级场景：恢复 preinst 保存的状态
if [ "$1" = "configure" ] && [ -n "$2" ] && [ -f "$STATE_FILE" ]; then
    echo "恢复 $SERVICE 升级前的状态..."
    
    # 从临时文件读取之前保存的状态
    . "$STATE_FILE"  # 加载 OLD_ENABLED 和 OLD_ACTIVE 变量
    
    # 恢复启用状态
    if [ "$OLD_ENABLED" = "enabled" ]; then
        systemctl enable "$SERVICE" >/dev/null 2>&1 || true
        echo "恢复启用状态: enabled"
    elif [ "$OLD_ENABLED" = "disabled" ]; then
        systemctl disable "$SERVICE" >/dev/null 2>&1 || true
        echo "恢复启用状态: disabled"
    else
        # 其他状态（如 masked、static 等）不主动修改，保持系统默认处理
        echo "保留原有特殊状态: $OLD_ENABLED"
    fi
    
    # 恢复运行状态
    if [ "$OLD_ACTIVE" = "active" ]; then
        systemctl restart "$SERVICE" >/dev/null 2>&1 || true
        echo "服务已重启（原状态为运行中）"
    else
        # 原状态为非运行中（如 inactive、failed），不主动启动
        echo "保留停止状态（原状态为非运行中）"
    fi
    
    # 清理临时状态文件（避免残留）
    rm -f "$STATE_FILE"
else
    # 处理首次安装场景：默认启用并启动服务（可根据需求调整）
    if [ "$1" = "configure" ] && [ -z "$2" ]; then
        echo "首次安装，初始化 $SERVICE 服务..."
        systemctl enable "$SERVICE" >/dev/null 2>&1 || true
        systemctl start "$SERVICE" >/dev/null 2>&1 || true
        echo "服务已默认启用并启动"
    fi
fi

exit 0
